# systemd/download-organizer.service
[Unit]
Description=Download Organizer — sort ~/Downloads by file type
Documentation=file://%h/Projects/download-organizer/README.md
# Only run when the user session is active (home directory is mounted).
After=default.target

[Service]
Type=oneshot

# Resolve the script path relative to this unit's location.
# %h expands to the user's home directory.
# Edit this path to match where you cloned the repository.
ExecStart=%h/Projects/download-organizer/organize-downloads.sh

# Stream stdout/stderr to journald with an identifier.
# Query with: journalctl --user -u download-organizer -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=download-organizer

# Harden the service: restrict write access to home only.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/Downloads %h/Projects/download-organizer/logs

# Do not restart on failure — this is a scheduled one-shot.
Restart=no

[Install]
WantedBy=default.target
